# systemd/wortschule.service
[Unit]
Description=wort.schule Rails Application
After=network.target

[Service]
Type=simple
User=wortschule
Group=wortschule
WorkingDirectory=/home/wortschule/app/current

# Load mise environment and start Rails
ExecStart=/bin/bash -lc 'eval "$($HOME/.local/bin/mise activate bash)" && mise trust --all && bundle exec puma -C config/puma.rb'

# Environment variables
Environment="RAILS_ENV=production"
Environment="MISE_YES=1"
EnvironmentFile=/home/wortschule/app/shared/.env

# Restart configuration
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/wortschule/app/shared/log/puma.log
StandardError=append:/home/wortschule/app/shared/log/puma.error.log

[Install]
WantedBy=multi-user.target